# airootfs/etc/systemd/system/reflector.service.d/override.conf
[Unit]
After=NetworkManager.service
Wants=NetworkManager.service
